# contrib/init/supermand.conf
description "Superman Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env SUPERMAND_BIN="/usr/bin/supermand"
env SUPERMAND_USER="superman"
env SUPERMAND_GROUP="superman"
env SUPERMAND_PIDDIR="/var/run/supermand"
# upstart can't handle variables constructed with other variables
env SUPERMAND_PIDFILE="/var/run/supermand/supermand.pid"
env SUPERMAND_CONFIGFILE="/etc/superman/superman.conf"
env SUPERMAND_DATADIR="/var/lib/supermand"

expect fork

respawn
respawn limit 5 120
kill timeout 600

pre-start script
    # this will catch non-existent config files
    # supermand will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$SUPERMAND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run supermand."
        echo "The setting must appear in $SUPERMAND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Superman Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$SUPERMAND_PIDDIR"
    chmod 0755 "$SUPERMAND_PIDDIR"
    chown $SUPERMAND_USER:$SUPERMAND_GROUP "$SUPERMAND_PIDDIR"
    chown $SUPERMAND_USER:$SUPERMAND_GROUP "$SUPERMAND_CONFIGFILE"
    chmod 0660 "$SUPERMAND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$SUPERMAND_PIDFILE" \
    --chuid $SUPERMAND_USER:$SUPERMAND_GROUP \
    --exec "$SUPERMAND_BIN" \
    -- \
    -pid="$SUPERMAND_PIDFILE" \
    -conf="$SUPERMAND_CONFIGFILE" \
    -datadir="$SUPERMAND_DATADIR" \
    -disablewallet \
    -daemon

